# systemd/grblwheel-kiosk.service
# GrblWheel kiosk: Chromium fullscreen to http://127.0.0.1:8765. Requires desktop (graphical.target).
[Unit]
Description=GrblWheel kiosk (Chromium fullscreen)
After=graphical.target grblwheel.service
Requires=grblwheel.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 3
ExecStart=/home/pi/GrblWheel/systemd/kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
